# stripSunrise.service
[Unit]
Description=Sends sunrise command to trinket

[Service]
Type=simple
ExecStart=/usr/local/bin/stripSunrise
