# systemd/stripSunset.service
[Unit]
Description=Sends sunset command to trinket

[Service]
Type=simple
ExecStart=/usr/local/bin/stripSunset
